# deployment/supervisor/address-validation-worker.conf
[program:address-validation-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/address-validation/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --memory=1024
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/address-validation/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=60

; Optional: Separate worker for exports (higher memory limit for large batches)
; Uncomment if exports are timing out or running out of memory
;
; [program:address-validation-export-worker]
; process_name=%(program_name)s_%(process_num)02d
; command=/usr/bin/php /var/www/address-validation/artisan queue:work --queue=exports --sleep=3 --tries=3 --max-time=3600 --memory=2048
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; user=www-data
; numprocs=1
; redirect_stderr=true
; stdout_logfile=/var/www/address-validation/storage/logs/export-worker.log
; stdout_logfile_maxbytes=10MB
; stopwaitsecs=120
